# .config/systemd/user/msmtp-queue.service
[Unit]
Description=Queued mail delivery service

[Service]
Type=oneshot
ExecStart=%h/.bin/msmtp-runqueue.sh
